# webos/deploy/prism-media-perms.conf
# /etc/event.d/prism-media-perms
# Grant world-rw on Qualcomm media/PMEM devices so apps outside the
# luna group can use OMX hardware video decode (libOmxCore.so).
# Install: novacom run file:///bin/cp -- /media/internal/prism-media-perms /etc/event.d/prism-media-perms

start on stopped udevd

script
    chmod a+rw /dev/pmem_adsp   2>/dev/null || true
    chmod a+rw /dev/msm_vidc_dec 2>/dev/null || true
end script
